# airootfs/etc/profile.d/start.sh
#!/bin/bash
ln -f /etc/os-release /usr/lib/os-release
clear
startx
